# whitelist.service
[Unit]

Description=Apply domain whitelist with ipset + iptables

After=network-online.target

Wants=network-online.target

 

[Service]

Type=oneshot

ExecStart=/usr/local/bin/whitelist.sh start

ExecStop=/usr/local/bin/whitelist.sh stop

RemainAfterExit=yes

 

[Install]

WantedBy=multi-user.target